# root/etc/s6-overlay/s6-rc.d/joal-fixattrs/run
#!/command/with-contenv sh

chown -R appuser:appuser /config
chmod -R g+w /config
