# root/etc/s6-overlay/s6-rc.d/joal-config/up
/etc/s6-overlay/s6-rc.d/joal-config/run
